# root/etc/services.d/mysql/run
#!/usr/bin/with-contenv sh

exec mysqld_safe
